# server/confs/super_uwsgi_mock.ini
[program:super_uwsgi]
command=/opt/pyEnvMock/bin/uwsgi --ini /opt/pyEnvMock/mainiway-django-mock/server/confs/uwsgi.ini
dirctory=/opt/pyEnvMock/mainiway-django-mock
environment=PATH="/opt/pyEnvMock/bin"
stdout_logfile=/opt/pyEnvMock/mainiway-django-mock/logs/super/super_uwsgi.access.log
stderr_logfile=/opt/pyEnvMock/mainiway-django-mock/logs/super/super_uwsgi.error.log
startsecs=0
stopwaitsecs=0
stopsignal=QUIT
autostart=true
autorestart=true
